# supervisord.conf
[supervisord]
nodaemon=true

[program:pyload]
command=python /pyload/pyLoadCore.py --configdir=/pyload/config --daemon
